# deploy/looseends-app.service
[Unit]
Description=Loose Ends — Slack agent (Bolt, Socket Mode)
After=network-online.target looseends-mcp.service
Wants=network-online.target
Requires=looseends-mcp.service

[Service]
Type=simple
User=looseends
WorkingDirectory=/opt/looseends
EnvironmentFile=/opt/looseends/.env
ExecStart=/opt/looseends/.venv/bin/python -u -m src.app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Exactly ONE instance, ever. Two Socket Mode connections split events between them:
# half your nudges vanish and buttons fire twice. systemd already guarantees this —
# do not put this behind an autoscaling group.

[Install]
WantedBy=multi-user.target
